# crontab
14 14 14 2 fri ~/Desktop/nomor1/nomor1.sh
2-30/6 * * * * ~/Desktop/nomor5/nomor5.sh
